# deploy/kalshibot.service
# systemd service so the bot runs 24/7 and auto-restarts.
# Copy to /etc/systemd/system/kalshibot.service on your server,
# then edit the two paths marked CHANGE-ME.
[Unit]
Description=Kalshi Bot
After=network-online.target
Wants=network-online.target

[Service]
# CHANGE-ME: the folder where you put the bot on the server
WorkingDirectory=/home/ubuntu/kalshi-bot
# CHANGE-ME: same folder + /.venv/bin/python
ExecStart=/home/ubuntu/kalshi-bot/.venv/bin/python run.py run --execute
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/kalshi-bot/logs/bot.out
StandardError=append:/home/ubuntu/kalshi-bot/logs/bot.out

[Install]
WantedBy=multi-user.target
